# deploy/site/exsender.service
[Unit]
Description=exsender control plane (FastAPI)
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=exsender
Group=exsender
WorkingDirectory=/opt/exsender/web
EnvironmentFile=/opt/exsender/web/.env
ExecStart=/opt/exsender/web/.venv/bin/python main.py
Restart=on-failure
RestartSec=5

NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
